# dr/core-backup-gw.timer
# /etc/systemd/system/core-backup-gw.timer
[Unit]
Description=Nightly production-core backup (via S3 File Gateway)

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
