# dr/core-backup.timer
# /etc/systemd/system/core-backup.timer
[Unit]
Description=Nightly production-core backup

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target